# themes/todaair/root/etc/uci-defaults/luci-theme-oxygen~
#!/bin/sh
uci batch <<-EOF
	set luci.themes.oxygen=/luci-static/oxygen
        commit luci
EOF
	
