# themes/todaair/root/etc/uci-defaults/luci-theme-todaair
#!/bin/sh
uci batch <<-EOF
	set luci.themes.todaair=/luci-static/todaair
        commit luci
EOF
	
